# screeps-stats/supervisord.conf
[supervisord]
nodaemon=true

[program:screeps_exporter]
command=python /app/screeps_exporter.py
autostart=true
autorestart=true
redirect_stderr=true

[program:prometheus]
command=/prometheus/prometheus \
    --config.file=/etc/prometheus/prometheus.yml \
    --storage.tsdb.path=/prometheus-data \
    --web.listen-address=:9090
autostart=true
autorestart=true
redirect_stderr=true
